# pfiles/standby.supervisord.conf
[supervisord]
nodaemon = true
logfile = /var/log/postgresql/supervisord.log
pidfile = /var/run/supervisor/supervisord.pid
user = root

[program:postgres]
command = /usr/local/bin/postgres.sh
startsecs = 0
autorestart = false
startretries = 1
user = postgres